# fan.service
[Unit]
Description=My RPI Fan Service
[Service]
ExecStart=/opt/fan.sh
